# supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".
[unix_http_server]
file=%(here)s/.runtime/supervisor.sock

[supervisord]
logfile=%(here)s/.log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=%(here)s/.runtime/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
directory=%(here)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/.runtime/supervisor.sock

[program:odoo]
command=%(here)s/bin/python odoo/odoo.py -c .openerp_serverrc --proxy-mode --load=web,connector
environment=ODOO_CONNECTOR_CHANNELS="root:4"
autostart=true
autorestart=true
stdout_logfile=%(here)s/.log/odoo.log
stderr_logfile=%(here)s/.log/odoo-error.log